# docker/dist/conf/sandbox.conf
[program:heron-tracker]
command = /usr/local/bin/heron-tracker --rootpath ~/.herondata/repository/state/local

[program:heron-ui]
command = /usr/local/bin/heron-ui

[program:heron-api-server]
command = /usr/local/bin/heron-apiserver --base-template=local --cluster=local -Dheron.directory.home=/usr/local/heron

autorestart = false
